# notifications_server/supervisord.conf
[program:notifications_server]
command=/usr/bin/node server.js
directory=/nodejs/notifications_server
user=nobody
logfile_maxbytes = 10M
stdout_logfile_backups = 2
